# conf/crontab
* * * * * /home/pi/phone/autoconnect/check-and-connect-bt-pan.sh